# conf/supervisor/electronicspi.ini
[program:electronicspi]
directory=/root/electronicspi_django/mysite
command=/root/electronicspi_django/env/bin/gunicorn -w 1 mysite.wsgi -b 127.0.0.1:9090
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
